# resources/tango-db.conf
[program:tango-db]
command=bin/bash -c '
  exec -- /usr/local/bin/DataBaseds 2 -ORBendPoint giop:tcp::$ORB_PORT'
autorestart=true
startretries=5
priority=1
